# sysconf/clearsync.service.in
[Unit]
Description=System synchronization daemon
After=syslog.target network.target

[Service]
Type=forking
User=clearsync
ExecStart=@sbindir@/clearsyncd -c @sysconfdir@/clearsync.conf
Restart=on-abort

[Install]
WantedBy=multi-user.target
